# openstack-venvs/packages/debian/neutron/debian/erp-neutron-compute-node-services.postinst
#! /bin/sh

for name in openvswitch-agent
do
    systemctl enable erp-neutron-$name
done
