# openstack-venvs/packages/debian/neutron/debian/erp-neutron-services.postinst
#! /bin/sh

for name in dhcp-agent l3-agent metadata-agent server
do
    systemctl enable erp-neutron-$name
done
