# openstack-venvs/packages/debian/cinder/debian/erp-cinder-compute-node-services.postinst
#! /bin/sh

for name in volume
do
    systemctl enable erp-cinder-$name
done
