# openstack-venvs/packages/debian/cinder/debian/erp-cinder-services.postinst
#! /bin/sh

for name in api scheduler
do
    systemctl enable erp-cinder-$name
done
