# openstack-venvs/packages/debian/nova/debian/erp-nova-compute-node-services.postinst
#! /bin/sh

for name in compute
do
    systemctl enable erp-nova-$name
done
