# openstack-venvs/packages/debian/nova/debian/erp-nova-services.postinst
#! /bin/sh

for name in api cert conductor consoleauth scheduler
do
    systemctl enable erp-nova-$name
done
