# openstack-venvs/services/erp-nova-conductor.service
[Unit]
Description=OpenStack Nova Conductor Server
After=syslog.target network.target

[Service]
Type=notify
NotifyAccess=all
TimeoutStartSec=0
Restart=always
User=nova
ExecStart=/srv/nova/bin/nova-conductor --config-file /etc/nova/nova.conf --log-file /var/log/nova/conductor.log

[Install]
WantedBy=multi-user.target

