# openstack-venvs/packages/debian/heat/debian/erp-heat-services.postinst
#! /bin/sh

for name in api api-cfn engine
do
    systemctl enable erp-heat-$name
done
