# openstack-venvs/packages/debian/glance/debian/erp-glance-services.postinst
#! /bin/sh

for name in api registry
do
    systemctl enable erp-glance-$name
done
